# scripts/client/install/systemd/orochi-agent-meta-push.service.template
[Unit]
Description=SciTeX Orochi — agent metadata heartbeat push (every 30s)
Documentation=https://github.com/ywatanabe1989/scitex-agent-container
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/systemd/user/orochi.env
ExecStart=/bin/bash -lc '@AGENT_META@ --push'
StandardOutput=journal
StandardError=journal
TimeoutStartSec=30
