# scripts/client/install/systemd/orochi-agent-meta-push.timer.template
[Unit]
Description=SciTeX Orochi — trigger agent metadata heartbeat every 30s

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=orochi-agent-meta-push.service
Persistent=true

[Install]
WantedBy=timers.target
