# deployment/fleet/launchd/com.scitex.orochi.AGENT_NAME.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template: replace AGENT_NAME with actual agent name.
     Install: cp to ~/Library/LaunchAgents/ then launchctl load <path>
     Generate all: deployment/fleet/install-launchd.sh -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.scitex.orochi.AGENT_NAME</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>source ~/.venv/bin/activate 2>/dev/null; exec scitex-agent-container start ~/.dotfiles/src/.scitex/orochi/agents/AGENT_NAME/AGENT_NAME.yaml</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>NetworkState</key>
    <true/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>120</integer>
  <key>WorkingDirectory</key>
  <string>/Users/ywatanabe</string>
  <key>StandardOutPath</key>
  <string>/Users/ywatanabe/.scitex/agent-container/logs/AGENT_NAME/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ywatanabe/.scitex/agent-container/logs/AGENT_NAME/launchd.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/sbin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/ywatanabe</string>
    <key>LANG</key>
    <string>en_US.UTF-8</string>
  </dict>
</dict>
</plist>
