# scripts/client/install/launchd/com.scitex.orochi.agent-meta-push.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Canonical agent-metadata heartbeat pusher for macOS (mba). Fires every
  30 seconds and POSTs agent state to the Orochi hub so the dashboard
  stays live. Parallel to orochi-agent-meta-push.{service,timer} on
  Linux hosts.

  The StartInterval=30 is the launchd equivalent of the systemd
  OnUnitActiveSec=30s timer.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.scitex.orochi.agent-meta-push</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>@AGENT_META@ --push</string>
  </array>
  <key>StartInterval</key>
  <integer>30</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>@HOME@/.scitex/agent-container/logs/orochi-agent-meta-push.out.log</string>
  <key>StandardErrorPath</key>
  <string>@HOME@/.scitex/agent-container/logs/orochi-agent-meta-push.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>@HOME@</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/sbin:/usr/sbin:/sbin</string>
    <key>LANG</key>
    <string>en_US.UTF-8</string>
  </dict>
</dict>
</plist>
